# settings/btrfs-backup-scripts/btrfs-snapshot-home.service
[Unit]
Description=Make btrfs snapshots for /home

[Service]
ExecStart=/bin/bash ~/.capture-snapshot.sh
