# btrfs/btrfs-auto-backup.service
[Unit]
Description=Make btrfs snapshots

[Service]
ExecStart=/usr/bin/python3 /home/kr2/Desktop/toys/btrfs-backup-scripts/btrfs-auto-backup.py
