# btrfs-backup-scripts/btrfs-auto-backup.timer
[Unit]
Description=Execute btrfs snapshot daily

[Timer]
OnBootSec=10min
OnUnitActiveSec=6hour

[Install]
WantedBy=timers.target
